# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
# ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin root - $TERM
ExecStart=-/sbin/agetty --autologin stratos --noclear %I 38400 linux
