# airootfs/etc/xdg/autostart/custom.desktop
# Start GNOME in the live environment
[Desktop Entry]
Type=Application
Exec=gnome-session